# src/backend/backend-node.service.example
# Example systemd unit to keep the Node backend running on your server.
#
# 1. Copy your built site to the server (the `out` folder).
# 2. Adjust WorkingDirectory below to point at out/backend on your server.
# 3. Copy this file to /etc/systemd/system/backend-node.service (drop the .example).
# 4. Enable and start it:
#       sudo systemctl daemon-reload
#       sudo systemctl enable --now backend-node
#       sudo systemctl status backend-node
#
# Check the Node path with `which node` and update ExecStart if it isn't /usr/bin/node
# (nvm installs live elsewhere, e.g. /home/<user>/.nvm/versions/node/<v>/bin/node).

[Unit]
Description=Nibula Node backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/your-site/out/backend
ExecStart=/usr/bin/node _core/index.js
Environment=PORT=3000
Environment=HOST=127.0.0.1
Environment=APP_ENV=production
Restart=on-failure
RestartSec=2
User=www-data

[Install]
WantedBy=multi-user.target
